# images/alpine-tedge-s6/s6-rc.d/tedge-mapper-c8y/run
#!/command/with-contenv sh
exec 2>&1
exec tedge-mapper c8y
