# images/alpine-s6/s6-rc.d/c8y-configuration-plugin-log/run
#!/bin/sh
sed 's/^/\x1b[34mconfiguration |\x1b[0m /' --unbuffered