# images/alpine-s6/s6-rc.d/mosquitto-log/run
#!/bin/sh
sed 's/^/\x1b[90mmosquitto     |\x1b[0m /' --unbuffered
