# images/alpine-s6/s6-rc.d/tedge-agent-log/run
#!/bin/sh
sed 's/^/\x1b[35magent         |\x1b[0m /' --unbuffered
