# images/alpine-s6/s6-rc.d/c8y-firmware-plugin-log/run
#!/bin/sh
sed 's/^/\x1b[33mfirmware      |\x1b[0m /' --unbuffered
