# images/alpine-s6/s6-rc.d/c8y-log-plugin-log/run
#!/bin/sh
sed 's/^/\x1b[32mlog           |\x1b[0m /' --unbuffered
