# images/alpine-tedge-s6/s6-rc.d/tedge-mapper-c8y-log/run
#!/bin/sh
sed 's/^/\x1b[36mmapper-c8y    |\x1b[0m /' --unbuffered
